# roulette.service
[Unit]
Description=Roulette Tracker Prediction Server
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/roulette-tracker
ExecStart=/opt/roulette-tracker/venv/bin/uvicorn server:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
